# docker/worker/crontab/root
* * * * * worker SLEEP=0 bash /docker-entrypoint.sh
